# run/sensu_server/etc/services.d/sensu-server/run
#!/usr/bin/execlineb -P
/opt/sensu/bin/sensu-server -c /etc/sensu/config.json -d /etc/sensu -e /etc/sensu/extensions -v
