# run/flapjack/etc/services.d/flapjack-httpchecker/run
#!/usr/bin/execlineb -P
/opt/flapjack/bin/flapjack receiver httpchecker
