# run/sensu_server/etc/services.d/redis/run
#!/usr/bin/execlineb -P
/usr/bin/redis-server
