# build/flapjack/data/etc/services.d/flapjack/run
#!/usr/bin/execlineb -P
/opt/flapjack/bin/flapjack server start --no-daemonize
