# build/flapjack/data/etc/services.d/redis/run
#!/usr/bin/execlineb -P
/opt/flapjack/embedded/bin/redis-server /opt/flapjack/embedded/etc/redis/redis-flapjack.conf
